# deploy/systemd/sillybot.container
[Unit]
Description=Sillybot Discord bot
Wants=network-online.target
After=network-online.target

[Container]
# Replace OWNER and VERSION with the published GHCR image you want to run.
Image=ghcr.io/OWNER/sillybot:VERSION
ContainerName=sillybot

# Token file pre-loaded into this rootless user's Podman secret storage.
Secret=discord_token,target=/run/secrets/discord_token

# Protected host data directory bind-mounted into the container. Z labels it
# for SELinux; U maps ownership for the image's non-root sillybot user.
Volume=/srv/sillybot/data:/var/lib/sillybot:Z,U

Environment=DISCORD_TOKEN_FILE=/run/secrets/discord_token
Environment=DATABASE_PATH=/var/lib/sillybot/sillybot.db
# Set true only once completed snapshots are exported and restore-tested off-host.
Environment=BACKUP_SNAPSHOTS_ENABLED=false
Environment=RUST_LOG=sillybot=info,serenity=warn
# Uncomment to fast-sync slash commands to one guild instead of globally.
# Environment=DEV_GUILD_ID=123456789012345678

# Match compose stop_grace_period so SIGTERM has time to close the Gateway.
StopTimeout=30

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
