# deploy/systemd/sillybot.container
[Unit]
Description=Sillybot Discord bot
Wants=network-online.target
After=network-online.target

[Container]
# Replace OWNER and VERSION with the published GHCR image you want to run.
Image=ghcr.io/OWNER/sillybot:VERSION
ContainerName=sillybot

# Token file pre-loaded with: podman secret create discord_token /etc/sillybot/discord_token
Secret=discord_token,target=/run/secrets/discord_token

# Protected host data directory bind-mounted into the container.
# Create with: install -d -o <svc-uid> -g <svc-gid> -m 0750 /srv/sillybot/data
Volume=/srv/sillybot/data:/var/lib/sillybot:Z

Environment=DISCORD_TOKEN_FILE=/run/secrets/discord_token
Environment=DATABASE_PATH=/var/lib/sillybot/sillybot.db
Environment=BACKUP_SNAPSHOTS_ENABLED=true
Environment=RUST_LOG=sillybot=info,serenity=warn
# Uncomment to fast-sync slash commands to one guild instead of globally.
# Environment=DEV_GUILD_ID=123456789012345678

# Match compose stop_grace_period so SIGTERM has time to close the Gateway.
StopTimeout=30

# Pull a newer digest when `podman auto-update` runs (set up the timer separately).
AutoUpdate=registry

[Service]
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=default.target
